# infra/ollama.service
[Unit]
Description=Ollama Service - Invoice Extractor OCR
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Nota: El usuario se reemplazará dinámicamente en setup.sh
User=alex
Group=alex
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=3
Environment="OLLAMA_HOST=0.0.0.0:11434"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama

# Limites de recursos (ajustar según necesidad)
MemoryLimit=8G
CPUQuota=200%

[Install]
WantedBy=multi-user.target

